# systemd/escucha.service
[Unit]
Description=Escucha speech-to-text service
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/escucha
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
